# ansible/roles/clicker/templates/clicker-backend.service.j2
[Unit]
Description=Clicker Game Backend
After=network.target

[Service]
Type=simple
WorkingDirectory={{ backend_dir }}
ExecStart=/usr/bin/node server.js
Restart=on-failure
User=www-data
Group=www-data
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=DB_DIR=/var/lib/clicker
# На всякий случай PATH, если node не /usr/bin/node
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target